# rsyslox.service
[Unit]
Description=rsyslox — syslog viewer and REST API
Documentation=https://github.com/phil-bot/rsyslox
After=network.target mysql.service mariadb.service
Wants=network.target

[Service]
Type=simple
User=rsyslox
Group=rsyslox
ExecStart=/opt/rsyslox/rsyslox
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/rsyslox
PrivateTmp=true

# Override config path if needed:
# Environment=RSYSLOX_CONFIG=/etc/rsyslox/config.toml

[Install]
WantedBy=multi-user.target
